# plugin/shadowsocks/template/supervisor.conf
[program:shadowsocks-<%= account.username %>]
command=ssserver -c /etc/shadowsocks/<%= account.username %>.json
autorestart=true
user=nobody
