# plugin/supervisor/template/program.conf
[program:<%= program.program_name %>]
command = <%= program.command %>
process_name = <%= program.name %>
autostart = <%= program.autostart.toString() %>
autorestart = <%= program.autorestart %>
user = <%= account.username %>
redirect_stderr = true
stdout_logfile = /home/<%= account.username %>/supervisor-<%= program.name %>.log
<% if(program.directory) { %>
    directory = <%= program.directory %>
<% } %>
